# .docker/supervisor/ssr.conf
[program:ssr]
process_name=%(program_name)s_%(process_num)02d
command=node /var/www/html/bootstrap/ssr/ssr.mjs
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
